# ops/backup-verify.service
[Unit]
Description=Vérification hebdomadaire restauration backup PostgreSQL
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=oneshot
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/app
Environment="BACKUP_DIR=/backups"
ExecStart=/home/ubuntu/app/ops/restore-db.sh latest
StandardOutput=journal
StandardError=journal
# Sécurité : écriture limitée. NoNewPrivileges NON activé car le script
# utilise ``sudo -u postgres`` pour créer/drop la DB de test (peer auth
# locale Postgres).
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/home/ubuntu/app /backups

[Install]
WantedBy=multi-user.target
