# ops/ramasse-purge.service
[Unit]
Description=Purge des ramasses soft-deleted expirées (>7 jours)
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=oneshot
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/app
# Utilise le venv du service ferment (mêmes dépendances : dotenv, sqlalchemy…)
EnvironmentFile=/home/ubuntu/app/.env
ExecStart=/home/ubuntu/app/venv/bin/python /home/ubuntu/app/scripts/purge_expired_ramasses.py 7
StandardOutput=journal
StandardError=journal
# Sécurité : pas de nouveau privilège, pas d'écriture ailleurs que /tmp
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/home/ubuntu/app

[Install]
WantedBy=multi-user.target
